# com.amir.jarvis.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amir.jarvis</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/amir/dev/personal/jarvis/.venv/bin/python3</string>
        <string>/Users/amir/dev/personal/jarvis/clap_detector.py</string>
    </array>

    <!-- Start automatically when the user logs in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if the process exits or crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Logs written by the script go here; stdout/stderr from Python itself go below -->
    <key>StandardOutPath</key>
    <string>/Users/amir/Library/Logs/jarvis/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/amir/Library/Logs/jarvis/stderr.log</string>

    <!-- LaunchAgents don't inherit the user shell PATH, so set it explicitly -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
